# ops/systemd/immcad-backend-local.service
[Unit]
Description=IMMCAD Local Backend (uvicorn) for Cloudflare named tunnel origin
After=network-online.target
Wants=network-online.target
PartOf=immcad-cloudflare-origin-stack.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/IMMCAD
Environment=PYTHONUNBUFFERED=1
ExecStart=/bin/bash -lc 'set -a; source /home/ec2-user/IMMCAD/ops/runtime/.env.backend-origin; set +a; export ENABLE_CASE_SEARCH=true; exec /home/ec2-user/IMMCAD/scripts/venv_exec.sh uvicorn immcad_api.main:app --app-dir src --host 127.0.0.1 --port 8002'
Restart=always
RestartSec=3
TimeoutStartSec=90
TimeoutStopSec=30
NoNewPrivileges=yes

[Install]
WantedBy=immcad-cloudflare-origin-stack.target
